# .config/systemd/user/gpg-agent.service
[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon --enable-ssh-support --homedir %h/.gnupg

[Install]
WantedBy=multi-user.target

