# .config/autostart/flashfocus.desktop
[Desktop Entry]
Version=1.0
Name=Flash Focus
Exec=/usr/bin/flashfocus
Terminal=false
Type=Application
OnlyShowIn=sway;
