# .config/autostart/flashfocus.desktop
[Desktop Entry]
Exec=/usr/bin/flashfocus
Name=flashfocus
Type=Application
